# scripts/auto_trade_kr.service
[Unit]
Description=Auto Trade KR - Korean Stock Trading Bot
After=network.target mysql.service

[Service]
Type=simple
User=brighter87
WorkingDirectory=/home/brighter87/Projects/asset
ExecStart=/home/brighter87/Projects/asset/venv/bin/python3 auto_trade.py
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto_trade_kr

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
